# packages/mistral/rpm/mistral.service
[Unit]
Description=Mistral workflow service
After=network.target

[Service]
Type=simple
User=mistral
Group=mistral
Environment="DAEMON_ARGS=--config-file /etc/mistral/mistral.conf --log-file /var/log/mistral/mistral.log"
EnvironmentFile=-/etc/sysconfig/mistral
ExecStart=/usr/share/python/mistral/bin/mistral-server $DAEMON_ARGS
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
