# packages/st2mistral/rpm/mistral.service
[Unit]
Description=Mistral workflow service invocation wrapper
After=network.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/mistral
ExecStartPost=/usr/bin/systemctl start mistral-api
ExecStartPost=/usr/bin/systemctl start mistral-server
ExecStop=/bin/sh -c "systemctl stop mistral-api; systemctl stop mistral-server;"
ExecStart=/bin/sh -c true
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
