# sources/st2reactor/debian/st2reactor.st2rulesengine.upstart
description     "StackStorm st2rulesengine service"
author          "StackStorm Ops"

start on filesystem and net-device-up IFACE!=lo
stop on starting rc RUNLEVEL=[016]

setuid st2
setgid st2
respawn
respawn limit 2 5

umask 007
kill timeout 60

script
  NAME=st2rulesengine
  DEFAULT_ARGS="--config-file /etc/st2/st2.conf"

  # Read configuration variable file if it is present
  [ -r /etc/default/$NAME ] && . /etc/default/$NAME

  /usr/bin/$NAME ${DEFAULT_ARGS}
end script
